# services/autorun.service
[Unit]
Description=Autorun applications

[Service]
User=lucicat
ExecStart=/usr/scripts/autorun.sh

[Install]
WantedBy=default.target
